# modules/buildduty_tools/templates/firefoxinfrachangelog.cron.erb
* */3 * * * python3 /home/buildduty/buildduty-tools/repos/firefoxinfrachangelog/update-files.py 2>&1